# deployment/systemd/tms-backend.service
# Systemd service file for TMS Backend (FastAPI)
# Location: /etc/systemd/system/tms-backend.service
#
# Installation:
#   sudo cp deployment/systemd/tms-backend.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable tms-backend
#   sudo systemctl start tms-backend
#
# Management:
#   sudo systemctl status tms-backend
#   sudo systemctl restart tms-backend
#   sudo systemctl stop tms-backend
#   sudo journalctl -u tms-backend -f

[Unit]
Description=TMS Backend (FastAPI)
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=tmsapp
Group=tmsapp
WorkingDirectory=/home/tmsapp/tms-server

# Environment file (choose one based on environment)
# For production:
EnvironmentFile=/home/tmsapp/tms-server/.env.production
# For staging:
# EnvironmentFile=/home/tmsapp/tms-server/.env.staging

# Activate virtual environment and run uvicorn
ExecStart=/home/tmsapp/tms-server/venv/bin/uvicorn app.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 4 \
    --log-level info \
    --access-log \
    --proxy-headers \
    --forwarded-allow-ips='*'

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=350
StartLimitBurst=10

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/tmsapp/tms-server

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tms-backend

# Environment variables that can be overridden
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/tmsapp/tms-server"

[Install]
WantedBy=multi-user.target
